# root/etc/services.d/dnsmasq/run
#!/usr/bin/with-contenv bash

  exec \
        dnsmasq -k
